# gawati/installer/gawatiserver/scripts/01/gawatiserver.service
[Unit]
Description=Gawati API server

[Service]
User=${RUNAS_USER}
WorkingDirectory=${SERVER_HOME}
ExecStart=/usr/bin/node ./bin/www

[Install]
WantedBy=multi-user.target

